# ansible/upstart.conf
description "Prince du Marais website"
author  "Benoit Hirbec <benoit.hirbec@gmail.com>"

start on vagrant-mounted
# stop on runlevel [0156]

chdir /vagrant
exec python website.py
respawn
